# tank/systemd/sync-tank.service
[Unit]
Description=Sync Tank camera/control service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=@SYNC_TANK_USER@
WorkingDirectory=@SYNC_TANK_ROOT@
Environment=PYTHONPATH=@SYNC_TANK_ROOT@
ExecStart=@SYNC_TANK_ROOT@/scripts/run-dev.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
